# docker/web/supervisor/uwsgi.conf
[program:uwsgi]
command=uwsgi /opt/uwsgi.ini --pidfile=/tmp/project-master.pid
